# supervisor.conf/celery_board_beat.conf
; ==================================
; celery worker supervisor example
; ==================================
; the name of your supervisord program
[program:boardbeat]

; set full path to celery program if using virtualenv
command=/usr/local/bin/celery -A election_boards beat -l info --pidfile="/var/run/celery/celerybeat.pid" -s /var/run/celery/celerybeat-schedule

; the directory to your django project
directory=/usr/src/app/election-boards

; if supervisord is run as the root user, switch users to this unix user account before doing any processing.
user=user

; supervisor will start as many instances of this program as named by numprocs
numprocs=1

; put process stdout output in this file
stdout_logfile=/var/log/celery/beat.log

; put process stderr output in this file
stderr_logfile=/var/log/celery/beat.log

; if true, this program will start automatically when supervisord is started
autostart=true

; may be one of false, unexpected, or true. if false, the process will never be autorestarted. if unexpected, the process will be restart when the program exits with an exit code that is not one of the exit codes associated with this process’ configuration (see exitcodes). if true, the process will be unconditionally restarted when it exits, without regard to its exit code.
autorestart=true

; the total number of seconds which the program needs to stay running after a startup to consider the start successful.
startsecs=10

; if your broker is supervised, set its priority higher so it starts first
priority=999
